# debian/ept-cache.cron.weekly
#!/bin/sh

# Update the data and index files of Debtags
test -x /usr/bin/ept-cache && /usr/bin/ept-cache reindex --quiet
