# debian/ept-cache.cron.weekly
#!/bin/sh

# Update the data and index files of Debtags
if [ -x /usr/bin/ept-cache ]; then
        /usr/bin/ept-cache reindex --quiet
fi
